# maintenance/sbin/cron.daily/laser_download_sentinel2
#!/bin/sh
#
# download Sentinel-2 metadata from AWS
#

/usr/bin/python /home/institut/rawdata/maintenance/rawdata/sat/sentinel2/02_aws_get_metadata.py --days 7 --quiet

